# alert/cpu_monitor.service
[Unit]
Description=CPU Load Monitor

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cpu_monitor.sh
